# _config/systemd/user/batwarn.service
[Unit]
Description=Battery level check

[Service]
Type=oneshot
ExecStart=%h/.dot-files/swaybatwarn.sh
